# assets/config/supervisor/1-init.conf
[program:init]
command=sh /scripts/init.sh
stderr_logfile=/var/log/init.err.log
exitcodes=0
startsecs=0
priority=1
